# homework-10/roles/deploy_unit/files/watchlog.service
[Unit]
Description=My watchlog service

[Service]
EnvironmentFile=/etc/default/watchlog
ExecStart=/opt/watchlog.sh $WORD $LOG
